# images/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
user = root
autostart = true

[program:php5-fpm]
command = /usr/sbin/php5-fpm -F
user = root
autostart = true

[program:cron]
command = cron -f -L 15
user = root
autostart=true
autorestart=true

[program:gearmanworker]
command=php /var/www/xtc/sh/gearmanworker.php
user=root
autostart=true
autorestart=true
stdout_logfile=/var/www/logs/gearmanworker.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/var/www/logs/gearmanworker.log
stderr_logfile_maxbytes=1MB

